# backup.service
[Unit]
Description=Backup service using restic

[Service]
ExecStart=/usr/local/bin/restic-backup