# backup-prune.timer
[Unit]
Description=Backup prune service timer

[Timer]
OnCalendar=Tue *-*-* 04:44:44
Unit=backup-prune.service

[Install]
WantedBy=multi-user.target