# backup.timer
[Unit]
Description=Backup service timer

[Timer]
#OnBootSec=0min
OnCalendar=2:22:22
Unit=backup.service

[Install]
WantedBy=multi-user.target