# backup-prune.service
[Unit]
Description=Backup service to prune restic repo

[Service]
ExecStart=/usr/local/bin/restic-prune